# site/profile/files/etc/profile.d/hadoop.sh

export HADOOP_INCLUDE_DIR=/opt/hadoop/include
export HADOOP_LIB_DIR=/opt/hadoop/lib/native
export JAVA_LIB_DIR=/usr/lib/jvm/jre
